# s6-rc.d/ood-update/up
#!/bin/bash -ex

/opt/ood/ood-portal-generator/sbin/update_ood_portal --force --insecure
