# s6-rc.d/ood-update/up
foreground { /mkdir -p /etc/ood/dex }
foreground { /chown -R ondemand-dex:ondemand-dex /etc/ood/dex }
foreground { /usr/bin/chmod -R 600 /etc/ood/dex }
foreground { /usr/bin/chmod 700 /etc/ood/dex }

/opt/ood/ood-portal-generator/sbin/update_ood_portal --force
