# s6-rc.d/apache/run
#!/bin/bash -ex

/usr/sbin/httpd -DFOREGROUND
