# s6-rc.d/apache/run
#!/bin/sh -ex

/usr/sbin/httpd -DFOREGROUND
